# crontab
0 7 * * * sh /mnt/stream.sh > /dev/stdout