# crontab
0 7 * * * sh /stream/stream.sh > /dev/stdout